# share/systemd/ascii-chat-server.service
[Unit]
Description=ascii chat Video Server
Documentation=man:ascii-chat(1)
Documentation=https://github.com/zfogg/ascii-chat
Documentation=https://zfogg.github.io/ascii-chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ascii-chat server
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening (user service compatible)
NoNewPrivileges=yes
ProtectHome=read-only
PrivateTmp=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ascii-chat-server

# Resource limits
LimitNOFILE=4096
LimitNPROC=64

[Install]
WantedBy=default.target
